# motd/update-motd.d/00-header
#!/bin/sh
COLUMNS=$(tput cols) 
cat << 'EOF'

  __  __                       _                   _ 
 |  \/  |                     | |                 | |
 | \  / | ___ _ __ ___   ___  | |     ___  _ __ __| |
 | |\/| |/ _ \ '_ ' _ \ / _ \ | |    / _ \| '__/ _' |
 | |  | |  __/ | | | | |  __/ | |___| (_) | | | (_| |
 |_|  |_|\___|_| |_| |_|\___| |______\___/|_|  \__,_|

EOF
