# motd/update-motd.d/2-fortune
#!/bin/sh
COLUMNS=$(tput cols) 
title=$(fortune|cowsay)
cat << EOF
`printf "%*s\n" $(((${#title}+$COLUMNS)/2)) "$title"`

EOF
